# CodeBehind/crontab
*/5 * * * * /usr/bin/python /home/pi/Documents/IOT2FinalProjectOfficial/IOTFinalProjectGroup8_Pi/CodeBehindgetReadings.py