# CodeBehind/crontab
*/5 * * * * /usr/bin/python /home/pi/Documents/IOTFinalProject/DHT11Test.py